# com.user.todoist-ai-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.todoist-ai-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/viktor_svirskyi/Projects/todoist-ai-agent/dist/index.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/viktor_svirskyi/Projects/todoist-ai-agent</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/viktor_svirskyi/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/viktor_svirskyi</string>
    <key>TODOIST_API_TOKEN</key>
    <string>9fd8fadbf12ceaf65a7deca707ec3665875fe3cb</string>
    <key>TODOIST_WEBHOOK_SECRET</key>
    <string>136e3c2b4f684b098896834aef4ce718</string>
    <key>BRAVE_API_KEY</key>
    <string>BSAHfy-MM9amXWj5yon6ipbSSixcEIO</string>
    <key>CLAUDE_BASE_URL</key>
    <string>http://127.0.0.1:8317/v1</string>
    <key>CLAUDE_API_KEY</key>
    <string>sk-local-proxy</string>
    <key>CLAUDE_MODEL</key>
    <string>claude-opus-4-6</string>
    <key>PORT</key>
    <string>9000</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/viktor_svirskyi/Library/Logs/todoist-ai-agent.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/viktor_svirskyi/Library/Logs/todoist-ai-agent.log</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
